# project-dev.service
[Unit]
Description=Project Development Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/start-https-server.sh
Restart=on-failure
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target 